# systemd/bbs-watcher.path
[Path]
PathModified=%h/.config/bbs

[Install]
WantedBy=default.target
