# systemd/bbs-watcher.path
[Path]
PathModified=%h/.config/bbs/bbs.json

[Install]
WantedBy=default.target
