# systemd/bbs-watcher.service
[Unit]
Description=bbs restarter

[Service]
Type=oneshot
ExecStart=systemctl --user reload bbs.service

[Install]
WantedBy=default.target
